# kernel/morse-bundle/files/etc/init.d/morsechipreset
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2023 Morse Micro Pty Ltd. All rights reserved.
#

START=09


boot() {
	/morse/scripts/chipreset.sh > /dev/kmsg 2>&1
}
