# netifd-morse/etc/uci-defaults/99-23.05-reconf
#!/bin/sh

. /lib/functions.sh

# netifd in OpneWrt 23 tries to reconf by default
# and doesn't set the interface down.
# Disable reconf for old configurations, so we can
# reload_configs properly

force_no_reconf() {
	local config="$1"
	config_get type "$config" type
	[ "$type" = "morse" ] || return 1

	uci set wireless."${config}".reconf="0"
}

config_load wireless

config_foreach force_no_reconf wifi-device